# scripts/systemd/terezia-grafana.service
# Legacy reference only. Do not install this unit as active project runtime.
[Unit]
Description=Chat Pref Grafana Service (legacy reference)
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=__PROJECT_DIR__
ExecStart=/usr/bin/docker compose -f __PROJECT_DIR__/docker-compose.yml up grafana
ExecStop=/usr/bin/docker compose -f __PROJECT_DIR__/docker-compose.yml down
Restart=always
RestartSec=10s
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=chat-pref-grafana-legacy

[Install]
WantedBy=multi-user.target
