# scripts/systemd/terezia-api.service
# Legacy reference only. Do not install this unit as active project runtime.
[Unit]
Description=Chat Pref API Service (legacy reference)
After=network.target postgresql.service docker.service
Requires=postgresql.service docker.service

[Service]
Type=simple
User=root
WorkingDirectory=__PROJECT_DIR__
Environment="PATH=__VENV_BIN__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=__VENV_BIN__/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5s
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=chat-pref-api-legacy

[Install]
WantedBy=multi-user.target
